# streamio_notifications/rootfs/etc/services.d/streamio_notifications/run
#!/usr/bin/with-bashio
bashio::log.info "Starting Stremio Notifications addon..."
exec python3 /app/main.py
